# apps/moodle-elearning/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=docker-php-fpm

[program:nginx]
command=nginx -g "daemon off;"
